# health-dashboard.service
[Unit]
Description=WeatherPi Health Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=cagdas
WorkingDirectory=/home/cagdas
ExecStart=/usr/bin/python3 /home/cagdas/health_dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target