# weatherkiosk-webserver.service
[Unit]
Description=Weather Kiosk Web Server
After=network.target

[Service]
Type=simple
User=cagdas
Group=cagdas
WorkingDirectory=/home/cagdas/weatherpi
ExecStart=/home/cagdas/weatherpi/venv/bin/python3 /home/cagdas/weatherpi/webserver.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target