# weatherpi-kiosk-robust.service
[Unit]
Description=WeatherPi Kiosk Display (Robust)
After=multi-user.target graphical.target lightdm.service nginx.service network-online.target
Wants=graphical.target network-online.target
Requires=nginx.service
StartLimitIntervalSec=0

[Service]
Type=exec
User=cagdas
Group=cagdas
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/cagdas/.Xauthority
WorkingDirectory=/home/cagdas

# Use our robust startup script
ExecStart=/home/cagdas/weatherpi_robust_startup.sh

# Robust restart policy
Restart=always
RestartSec=30
StartLimitBurst=10

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal

# Security (but allow X11 access)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/cagdas /var/log

[Install]
WantedBy=graphical.target