# pi_configs/weatherpi-kiosk.service
[Unit]
Description=WeatherPi Kiosk Display
After=graphical.target lightdm.service
Wants=graphical.target

[Service]
Type=simple
# Run as the user's account on this device
User=cagdas
Environment=DISPLAY=:0
WorkingDirectory=/home/cagdas
ExecStartPre=/bin/sleep 15
# By default use Firefox ESR - on Pi4 you may prefer the Chromium unit provided in pi_configs/
ExecStart=/bin/bash -c "firefox-esr --kiosk http://localhost/weather.html"
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target
